# infrastructure/config/systemd/rateguru-mailtrap-local.service
[Unit]
Description=RateGuru staging Mailtrap Local (secondary mail-capture mirror)
Documentation=https://github.com/mailtrap/mailtrap-local
Documentation=file:/home/www/rateguru/runbooks/mail-capture.md
# This service is independent: it must never depend on Mailpit. Mailpit
# relays a best-effort mirror to it, not the other way around.
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=rateguru-mailtrap-local
Group=rateguru-mailtrap-local
# Loopback-only listeners and an explicit SQLite path inside the state dir.
ExecStart=/usr/local/bin/rateguru-mailtrap-local \
    --smtp-listen 127.0.0.1:3535 \
    --http-listen 127.0.0.1:3550 \
    --db /var/lib/rateguru-mail-capture/mailtrap-local/mailtrap-local.sqlite3
# Non-secret storage config (max_messages) and per-user state locations.
Environment=HOME=/var/lib/rateguru-mail-capture/mailtrap-local
Environment=XDG_CONFIG_HOME=/var/lib/rateguru-mail-capture/mailtrap-local
Environment=MAILTRAP_LOCAL_CONFIG=/etc/rateguru/mail-capture/mailtrap-local.yml
Environment=MAILTRAP_LOCAL_SECRET_KEY_FILE=/var/lib/rateguru-mail-capture/mailtrap-local/secret.key
WorkingDirectory=/var/lib/rateguru-mail-capture/mailtrap-local
Restart=on-failure
RestartSec=2

# Logs go to journald (Mailtrap Local logs to stderr).
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rateguru-mailtrap-local

# --- Hardening -------------------------------------------------------------
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
# Loopback TCP only (plus AF_UNIX for internal use); no raw/netlink families.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
# Empty capability set: Mailtrap Local binds only high ports (>1024).
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0077
# The persistent SQLite directory is the only writable location.
ReadWritePaths=/var/lib/rateguru-mail-capture/mailtrap-local

[Install]
WantedBy=multi-user.target
